# kiosk.service
[Unit]
Description=Touch pHAT Video Kiosk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi

WorkingDirectory=/home/pi/kiosk

# Pull latest app + config from git. Non-fatal so the kiosk still
# starts on stale files when there's no network.
ExecStartPre=/bin/sh -c 'git -C /home/pi/kiosk pull || true'

# Update yt-dlp before starting the kiosk.
# Failures are logged but do not prevent startup.
ExecStartPre=/bin/sh -c '/usr/local/bin/yt-dlp -U 2>&1 | /usr/bin/logger -t kiosk-yt-dlp || true'

# Run the kiosk using the touchphat venv Python
ExecStart=/home/pi/touchphat/venv/bin/python /home/pi/kiosk/kiosk.py

# If the kiosk crashes, restart it after 10 seconds
Restart=on-failure
RestartSec=10

# DRM framebuffer access
SupplementaryGroups=video render tty

Environment=HOME=/home/pi
Environment=PYTHONUNBUFFERED=1

# stdout is connected directly to tty2 by systemd (as root, before dropping
# to the pi user) so the marquee renders on the console. stderr still goes
# to the journal, keeping logs and display cleanly separated.
StandardOutput=tty
TTYPath=/dev/tty2
StandardError=journal
TTYReset=yes
TTYVTDisallocate=yes
SyslogIdentifier=kiosk

[Install]
WantedBy=multi-user.target
